# curator/image/etc/cron.d/curator
20 23 * * * root /usr/local/bin/curator --config /etc/curator/curator.yml /etc/curator/actions.yml
